# udev/99-drm-root.rules
SUBSYSTEM=="drm", KERNEL=="card1", RUN+="/usr/bin/setfacl -m u:root:rw /dev/%k"
